# deb-pack/workspaces/DEBIAN/prerm
#!/bin/bash

systemctl disable caddy
systemctl stop caddy